# debian/cron.d/debci-worker
21 */6 * * * root [ -x /usr/share/debci/bin/debci-setup ] && /usr/share/debci/bin/debci-setup

# vim: ft=crontab
