# debian/cron.d/debci-worker
21 */6 * * * debci [ -x /usr/share/debci/bin/debci-setup ] && sudo /usr/share/debci/bin/debci-setup

# vim: ft=crontab
