# debian/debci-generate-index.timer
[Unit]
Description=debci index generator timer

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
