# debian/debci-collector.cron.d
*/5 * * * * debci [ -x /usr/bin/debci ] && /usr/bin/debci generate-index --quiet

# vim: ft=crontab
