# debian/debci-generate-index.service
[Unit]
Description=debci index generator

[Service]
Type=oneshot
User=debci
ExecStart=/bin/sleep 1s
ExecStart=/usr/bin/debci generate-index
