# debian/debci-update.service
[Unit]
Description=debci index generator

[Service]
Type=oneshot
User=debci
Group=debci
ExecStart=/bin/sleep 1s
ExecStart=/usr/bin/debci update
